# dotfiles/.config/conky/conky.service
[Unit]
Description=Conky Setup

[Service]
ExecStart=/home/csa/.config/conky/conky.sh

[Install]
WantedBy=multi-user.target
